# bonsAI_LLM/systemd/bonsai-llm.service
[Unit]
Description=bonsAI llama.cpp local LLM server
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/opt/bonsAI/bonsAI_LLM
Environment=BONSAI_LLM_ENV_FILE=/opt/bonsAI/bonsAI_LLM/.env
ExecStart=/opt/bonsAI/bonsAI_LLM/scripts/start_llm_server.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
